# systemd/triade-chat-ui.service
[Unit]
Description=Tríade Omega Local Chat UI
After=network.target triade-api.service
Wants=triade-api.service

[Service]
Type=simple
User=santiago
WorkingDirectory=/home/santiago/triadees
ExecStart=/home/santiago/triadees/.venv/bin/python -m uvicorn apps.chat_ui_app:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
