# systemd/triade-model-router.service
[Unit]
Description=Tríade Omega Model Router API (DEPRECATED — merged into single_port_app:8010)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=santiago
WorkingDirectory=/home/santiago/triadees
ExecStart=/home/santiago/triadees/.venv/bin/python -m uvicorn apps.single_port_app:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
